# update_server/deploy/amap-companion-sync.service
[Unit]
Description=Sync GitHub Release APK for AMap Companion
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/amap-companion/update_server
Environment=BRANCH=master
Environment=REPO_DIR=/opt/amap-companion
# Override this only if you fork the repository.
Environment=GITHUB_REPO=zuo-qirun/amap-companion
# Environment=GITHUB_TOKEN=ghp_xxx
ExecStart=/usr/bin/node sync-build.js
User=pi
